# cd/Implementierung/machines/app/contents/etc/supervisor/conf.d/php.conf
[program:php]
command=/usr/sbin/php5-fpm -F